# infra/livekit/host/picom-livekit-host-health.service
[Unit]
Description=Picom LiveKit host baseline health check
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/picom-livekit-host-health
User=root
Group=root
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/log

[Install]
WantedBy=multi-user.target
