# services/secret-sms-gateway/systemd/picom-secret-sms-gateway.service
[Unit]
Description=Picom self-hosted secret community SMS verification gateway
After=network-online.target kannel.service
Wants=network-online.target

[Service]
Type=simple
User=picom-sms
Group=picom-sms
EnvironmentFile=/etc/picom/secret-sms-gateway.env
ExecStart=/usr/bin/node /opt/picom/secret-sms-gateway/server.mjs
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
UMask=0077

[Install]
WantedBy=multi-user.target