# services/secret-voice-gateway/systemd/picom-secret-voice-gateway.service
[Unit]
Description=Picom self-hosted secret community voice verification gateway
After=network-online.target asterisk.service
Wants=network-online.target
Requires=asterisk.service
[Service]
Type=simple
User=picom-verify
Group=picom-verify
EnvironmentFile=/etc/picom/secret-voice-gateway.env
ExecStart=/usr/bin/node /opt/picom/secret-voice-gateway/server.mjs
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
[Install]
WantedBy=multi-user.target